# deploy/marana-server.service
[Unit]
Description=Andor Marana camera server (ZeroMQ)
Documentation=file:///opt/marana/marana_console/README.md
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify: the server sends READY=1 once up and pings WATCHDOG=1 while the
# camera worker is healthy. If the worker wedges (e.g. an SDK/USB acquisition
# hang), pings stop and systemd restarts the service after WatchdogSec.
# Recovery ≈ stale_after (10s, in the server) + WatchdogSec ≈ 30s. WatchdogSec
# must stay > the server's ping interval (5s) with margin.
Type=notify
WatchdogSec=20
User=andor
Group=andor
WorkingDirectory=/opt/marana/marana_console
Environment="PATH=/opt/marana/venv/bin:/usr/bin"
Environment="EPICS_CA_AUTO_ADDR_LIST=YES"
# Optional explicit broadcast (uncomment + set for your subnet if auto-discovery
# isn't picking up the IOC):
#Environment="EPICS_CA_ADDR_LIST=192.168.1.255"
ExecStart=/opt/marana/venv/bin/python -m marana_server \
  --bind 0.0.0.0 \
  --ctrl-port 5555 \
  --frame-port 5556 \
  --captures-dir /data/marana_captures \
  --log-level INFO
Restart=on-failure
RestartSec=2
StandardOutput=append:/var/log/marana/server.log
StandardError=append:/var/log/marana/server.log

[Install]
WantedBy=multi-user.target
